# etc/systemd/system/getty@tty1.service.d/override.conf.in
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin M4_USER --noclear %I $TERM
