# etc/udev/rules.d/00-usb-keyboards.rules.in
SUBSYSTEM=="input", \
          ACTION=="add", \
          ATTRS{idProduct}=="0125", \
          ATTRS{idVendor}=="04d9", \
          RUN+="/bin/su M4_USER M4_HOME/bin/udev-setkeyboard"
SUBSYSTEM=="input", \
          ACTION=="remove", \
          ATTRS{idProduct}=="0125", \
          ATTRS{idVendor}=="04d9", \
          RUN+="/bin/su M4_USER M4_HOME/bin/udev-setkeyboard unlock"
